# RaspberryPi_MP3_Player_BuildRoot/overlay/etc/profile.d/postbuild.sh
#!/bin/sh
ifup -a &>/dev/null
hello.o
cd /myApplications
./USBDemon.sh &
#./sound_player.sh &
